# openrc/concierge
#!/sbin/openrc-run

supervisor=supervise-daemon

description="Concierge HTTP(S) server"
command="/usr/local/bin/concierge"
command_args="-c /etc/concierge/concierge.toml"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
	need net
	after firewall
}