# openrc/concierge
#!/sbin/openrc-run

supervisor=supervise-daemon

command="/usr/local/bin/concierge"
command_args="-f /var/log/concierge.log -a /ssh:localhost:22 -crt /root/server.crt -key /root/server.key"
description="Concierge HTTP(S) server"

depend() {
	need net
	after firewall
}